# mcp/registry/templates/supervisor.conf.j2
# Supervisor configuration template for MCP
# Auto-generated by Windsurf MCP Registry
# Template version: 1.0

[unix_http_server]
file=/var/run/supervisor.sock   ; path to your socket file
chmod=0700                       ; socket file mode (default 0700)

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log ; main log file
logfile_maxbytes=50MB            ; max main logfile bytes b4 rotation
logfile_backups=10               ; # of main logfile backups
loglevel=info                    ; log level; default info; others: debug,warn,trace
pidfile=/var/run/supervisord.pid ; supervisord pidfile
user=root                       ; default user

; The rpcinterface:supervisor section must remain in the config file for
; RPC (supervisorctl/web interface) to work, additional interfaces may be
; added by defining them in separate [rpcinterface:x] sections.
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; The supervisorctl section configures how supervisorctl will connect to
; supervisord.  configure it match the settings in either the unix_http_server
; or inet_http_server section.
[supervisorctl]
serverurl=unix:///var/run/supervisor.sock ; use a unix:// URL  for a unix socket

; The [include] section can just contain the "files" setting.  This
; setting can list multiple files (separated by whitespace or newlines).
; It can also contain wildcards.  The filenames are interpreted as
; relative to this file.  Included files *cannot* include files themselves.
[include]
files = /etc/supervisor/conf.d/*.conf

; Sample program configuration for an MCP service
[program:{{ mcp_name }}]
command={{ mcp_command }}  ; the program (relative uses PATH, can take args)
process_name=%(program_name)s ; process_name expr (default %(program_name)s)
directory=/var/run/{{ mcp_name }} ; directory to cwd to before exec (def no cwd)
user={{ mcp_user | default('root') }} ; setuid to this UNIX account
numprocs=1                    ; number of processes copies to start
autostart=true                ; start at supervisord start (default: true)
autorestart=unexpected        ; when to restart if exited after running (default: unexpected)
startsecs=10                  ; number of secs prog must stay running (def. 1)
startretries=3                ; max # of serial start failures (default 3)
stopsignal=TERM               ; signal used to kill process (default TERM)
stopwaitsecs=10               ; max num secs to wait b4 SIGKILL (default 10)
stdout_logfile=/var/log/{{ mcp_name }}/stdout.log ; stdout log path, NONE for none;
stdout_logfile_maxbytes=1MB   ; max # logfile bytes b4 rotation (default 50MB)
stdout_logfile_backups=10     ; # of stdout logfile backups (default 10)
stderr_logfile=/var/log/{{ mcp_name }}/stderr.log ; stderr log path, NONE for none;
stderr_logfile_maxbytes=1MB   ; max # logfile bytes b4 rotation (default 50MB)
stderr_logfile_backups=10     ; # of stderr logfile backups (default 10)
